# packages/runner/systemd/runner.timer
[Unit]
Description=xGov Committees Service Runner Timer

[Timer]
OnUnitInactiveSec=50min
OnBootSec=0

[Install]
WantedBy=timers.target
