# runner/systemd/runner.timer
[Unit]
Description=xGov Committee Service Runner Timer

[Timer]
OnUnitInactiveSec=50min
OnBootSec=1s

[Install]
WantedBy=timers.target
